# dot_config/systemd/user/khal-notify.service
[Unit]
Description=Sync vdirsyncer and check khal for upcoming events

[Service]
Type=oneshot
ExecStart=%h/.local/bin/khal-notify.sh
